# crontab
#!/bin/sh
cd /var/lib/auth_token/app && update_geoip_db.sh
